# .github/scripts/auto-commit.service
[Unit]
Description=Auto-commit bootcamp changes every 5 minutes
Documentation=https://github.com/user/repo

[Service]
Type=oneshot
User=epti
WorkingDirectory=/home/epti/Documentos/sena/4T25/oop_java
ExecStart=/bin/bash /home/epti/Documentos/sena/4T25/oop_java/.github/scripts/auto-commit.sh
StandardOutput=append:/home/epti/Documentos/sena/4T25/oop_java/.github/scripts/auto-commit.log
StandardError=append:/home/epti/Documentos/sena/4T25/oop_java/.github/scripts/auto-commit.log

[Install]
WantedBy=default.target
